# dist/systemd/vector.service
[Unit]
Description=Vector
Documentation=http://vectorproject.io
After=network-online.target
Requires=network-online.target

[Service]
User=vector
Group=vector
ExecStart=/opt/vector/bin/vector --config=/etc/vector/vector.toml
Restart=no
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=vector

[Install]
WantedBy=multi-user.target
